# misc/systemd-service/dde-cooperation.service.in
[Unit]
Description=DDE cooperation daemon
After=network.target

[Service]
Type=simple
ExecStart=@libexecdir@/dde-cooperation/dde-cooperation-daemon

[Install]
WantedBy=multi-user.target
Alias=dbus-com.deepin.Cooperation.service
