# misc/systemd-service/dde-cooperation.service.in
[Unit]
Description=DDE cooperation
After=network.target

[Service]
Type=simple
ExecStart=@executable_install_dir@/dde-cooperation

[Install]
WantedBy=multi-user.target
Alias=dbus-org.deepin.dde.Cooperation1.service
