# src/apps/daemon/service/cooperation-daemon.service
[Unit]
Description=Deepin Cooperation Daemon

[Service]
ExecStart=/usr/bin/cooperation-daemon

[Install]
WantedBy=multi-user.target
